# hermes-zalo.service
[Unit]
Description=Hermes-Zalo - Zalo <-> Hermes Gateway
After=network.target mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hermes-zalo
EnvironmentFile=/opt/hermes-zalo/.env
ExecStart=/opt/hermes-zalo/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
